# crontab
0 6 * * * sh /var/scripts/db_backup.sh >/dev/null 2>&1
